# gunicorn.service
[Unit]
Description=gunicorn daemon for Django HRMS
After=network.target mysql.service redis.service
Requires=mysql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/var/www/hrms
Environment="PATH=/var/www/hrms/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=hrms.settings"
ExecStart=/var/www/hrms/venv/bin/gunicorn \
    --config gunicorn.py \
    --capture-output \
    --enable-stdio-inheritance \
    hrms.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/hrms /var/log/gunicorn

[Install]
WantedBy=multi-user.target